# feeds/net/mpcs/files/mpcs.init
#!/bin/sh /etc/rc.common

START=50

start () {
	[ "`cat /etc/inittab | grep #tty`" ] || sed -i 's/tty/#tty/g' /etc/inittab
	service_start /usr/sbin/mpcs -c /etc/mpcs &
}

stop () {
	service_stop /usr/sbin/mpcs
}
